# crontab
0 0 * * * /bin/bash /opt/check_raid.sh >/dev/null 2>&1